# etc/supervisor/conf.d/utorrent-server.sv.conf
[program:utorrent-server]
command=cd /opt/utorrent-server && /opt/utorrent-server/utserver -settingspath=/opt/utorrent-server/
redirect_stderr=true
