# etc/supervisor/conf.d/utorrent-server.sv.conf
[program:utorrent-server]
directory=/opt/utorrent-server
command=/opt/utorrent-server/utserver -settingspath /opt/utorrent-server/ -configfile /opt/utorrent-server/utserver.conf
redirect_stderr=true
