# etc/supervisor/conf.d/utserver.sv.conf
[program:utserver]
directory=/opt/utserver
command=/opt/utserver/utserver -settingspath /mnt/data/settings -logfile /mnt/data/log/utserver.log
redirect_stderr=true
