# etc/supervisor/conf.d/utserver.sv.conf
[program:utserver]
directory=/opt/utserver
command=/opt/utserver/utserver -settingspath /opt/utserver -logfile /opt/utserver/utserver.log
redirect_stderr=true
